# deploy/launchd/ai.csagent.backup-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.csagent.backup-weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__CSAGENT_ROOT__/deploy/backup-personal.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__CSAGENT_ROOT__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>CSAGENT_HOME</key>
        <string>__CSAGENT_HOME__</string>
        <key>CSAGENT_ROOT</key>
        <string>__CSAGENT_ROOT__</string>
        <key>PATH</key>
        <string>__PATH_SNAPSHOT__</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__CSAGENT_LOG_DIR__/backup-weekly.log</string>

    <key>StandardErrorPath</key>
    <string>__CSAGENT_LOG_DIR__/backup-weekly.error.log</string>
</dict>
</plist>
