# deploy/launchd/ai.irida.cron-tick.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.irida.cron-tick</string>

    <key>ProgramArguments</key>
    <array>
        <string>__NODE_BIN__</string>
        <string>__IRIDA_ROOT__/node_modules/.bin/tsx</string>
        <string>__IRIDA_ROOT__/src/cli.ts</string>
        <string>cron</string>
        <string>tick</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__IRIDA_ROOT__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>IRIDA_HOME</key>
        <string>__IRIDA_HOME__</string>
        <key>IRIDA_ROOT</key>
        <string>__IRIDA_ROOT__</string>
        <key>IRIDA_DATABASE_URL</key>
        <string>__IRIDA_DATABASE_URL__</string>
        <key>IRIDA_SECRETS_KEY</key>
        <string>__IRIDA_SECRETS_KEY__</string>
        <key>OBSIDIAN_VAULT_PATH</key>
        <string>__OBSIDIAN_VAULT_PATH__</string>
        <key>PATH</key>
        <string>__PATH_SNAPSHOT__</string>
    </dict>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__IRIDA_LOG_DIR__/cron-tick.log</string>

    <key>StandardErrorPath</key>
    <string>__IRIDA_LOG_DIR__/cron-tick.error.log</string>
</dict>
</plist>
