# deploy/launchd/ai.irida.digest-qa-morning.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.irida.digest-qa-morning</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__IRIDA_ROOT__/deploy/digest-qa-morning.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__IRIDA_ROOT__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>IRIDA_HOME</key>
        <string>__IRIDA_HOME__</string>
        <key>IRIDA_ROOT</key>
        <string>__IRIDA_ROOT__</string>
        <key>PATH</key>
        <string>__PATH_SNAPSHOT__</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__IRIDA_LOG_DIR__/digest-qa-morning.log</string>

    <key>StandardErrorPath</key>
    <string>__IRIDA_LOG_DIR__/digest-qa-morning.error.log</string>
</dict>
</plist>
